# ctrl/repo-server.conf
[program:nginx]
command=/usr/sbin/nginx
sutostart=true
autorestart=true
